# ansible/roles/beszel-agent/templates/beszel-agent.service.j2
[Unit]
Description=Beszel Agent Service
Wants=network-online.target
After=network-online.target

[Service]
EnvironmentFile=/etc/default/beszel-agent
ExecStart={{ beszel_agent_install_dir }}/beszel-agent
User={{ beszel_agent_user }}
Restart=on-failure
RestartSec=5
StateDirectory=beszel-agent

# Sandboxing (close to upstream install-agent.sh, plus PrivateTmp)
# PrivateTmp gives the unit its own /tmp namespace - required when /tmp is
# part of the root filesystem (no separate /tmp mount), since
# ProtectSystem=strict otherwise renders /tmp read-only and the agent
# fails on `/tmp/beszel_health` cleanup.
KeyringMode=private
LockPersonality=yes
PrivateTmp=true
ProtectClock=yes
ProtectHome=read-only
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectSystem=strict
RemoveIPC=yes
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
