# ansible/roles/ts-static-serves/files/ts-hubs.service
[Unit]
Description=Regenerate Tailscale hub (mini-homepage) HTML
After=tailscaled.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/ts-gen-hubs
